# scripts/patronum.service
[Unit]
Description=Patronum AI Agent
After=network.target

[Service]
Type=simple
User=patronum
WorkingDirectory=/var/lib/patronum
ExecStart=/usr/bin/node /var/lib/patronum/source/dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
